# Log Collect Automation/404L/404L 1st #7 log collect scheduler.xml
<?xml version="1.0" encoding="UTF-16"?> 
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task"> 
	<RegistrationInfo> 
		<Date>2023-03-08T20:42:39.6421403</Date> 
		<Author>RADIANT VISION SYSTEMS</Author> 
		<URI>Demura Log collect</URI> 
	</RegistrationInfo> 
	<Triggers> 
		<CalendarTrigger> 
			<StartBoundary>2023-03-08T02:30:00</StartBoundary> 
			<Enabled>true</Enabled> 
			<ScheduleByDay> 
				<DaysInterval>1</DaysInterval> 
			</ScheduleByDay> 
		</CalendarTrigger> 
	</Triggers> 
	<Principals> 
		<Principal id="Author"> 
			<UserId>S-1-5-21-725480862-3868392371-2725245176-1001</UserId> 
			<LogonType>InteractiveToken</LogonType> 
			<RunLevel>HighestAvailable</RunLevel> 
		</Principal> 
	</Principals> 
	<Settings> 
		<MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy> 
		<DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries> 
		<StopIfGoingOnBatteries>true</StopIfGoingOnBatteries> 
		<AllowHardTerminate>true</AllowHardTerminate> 
		<StartWhenAvailable>true</StartWhenAvailable> 
		<RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable> 
		<IdleSettings> 
			<StopOnIdleEnd>true</StopOnIdleEnd> 
			<RestartOnIdle>false</RestartOnIdle> 
		</IdleSettings> 
		<AllowStartOnDemand>true</AllowStartOnDemand> 
		<Enabled>true</Enabled> 
		<Hidden>false</Hidden> 
		<RunOnlyIfIdle>false</RunOnlyIfIdle> 
		<WakeToRun>false</WakeToRun> 
		<ExecutionTimeLimit>PT72H</ExecutionTimeLimit> 
		<Priority>7</Priority> 
	</Settings> 
	<Actions Context="Author"> 
		<Exec> 
			<Command>"D:\Program\RVS\Demura to Main 404L 1st #7.bat"</Command> 
		</Exec> 
	</Actions> 
</Task> 
